# supervisor.conf
[program:1jingdian]
command=/var/www/1jingdian/gunicorn_start
user=root
autostart=true
autorestart=true